# ansible/roles/caddy_l4/templates/channel-c-squid-ipv4-egress.service.j2
[Unit]
Description=Channel C Squid IPv4-only egress guard
Before=channel-c-squid.service
After=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -ec '/usr/sbin/ip6tables -C OUTPUT -m owner --uid-owner proxy -p tcp -j REJECT --reject-with tcp-reset 2>/dev/null || /usr/sbin/ip6tables -I OUTPUT 1 -m owner --uid-owner proxy -p tcp -j REJECT --reject-with tcp-reset'
ExecStop=/bin/sh -ec 'while /usr/sbin/ip6tables -C OUTPUT -m owner --uid-owner proxy -p tcp -j REJECT --reject-with tcp-reset 2>/dev/null; do /usr/sbin/ip6tables -D OUTPUT -m owner --uid-owner proxy -p tcp -j REJECT --reject-with tcp-reset; done; true'

[Install]
WantedBy=multi-user.target
