# ansible/roles/caddy_l4/templates/channel-c-stunnel.service.j2
[Unit]
Description=Channel C TLS wrapper for local forward proxy
{% if channel_c_naive_backend | default('stunnel_tinyproxy') == 'stunnel_squid' %}
After=network-online.target channel-c-squid.service
Wants=network-online.target channel-c-squid.service
{% else %}
After=network-online.target channel-c-tinyproxy.service
Wants=network-online.target channel-c-tinyproxy.service
{% endif %}

[Service]
Type=simple
ExecStart=/usr/bin/stunnel4 /etc/stunnel/channel-c.conf
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
