# ansible/roles/caddy_l4/templates/channel-c-squid.service.j2
[Unit]
Description=Channel C local Squid compatibility proxy
After=network-online.target
Wants=network-online.target
{% if channel_c_naive_squid_ipv4_only | default(true) | bool %}
After=channel-c-squid-ipv4-egress.service
Wants=channel-c-squid-ipv4-egress.service
{% endif %}

[Service]
Type=simple
ExecStart=/usr/sbin/squid -N -f /etc/squid/channel-c.conf
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
